# bin/com.tashian.autolights.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.tashian.autolights</string>
	<key>ProgramArguments</key>
	<array>
		<string>sh</string>
		<string>/usr/local/bin/autolights.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/autolights.job.err</string>
	<key>StandardOutPath</key>
	<string>/tmp/autolights.job.out</string>
</dict>
</plist>
